# Configs/autostart/.config/autostart/org.keepassxc.KeePassXC.desktop
[Desktop Entry]
Type=Application
Name=org.keepassxc.KeePassXC
Exec=flatpak run org.keepassxc.KeePassXC
X-Flatpak=org.keepassxc.KeePassXC
